# support/housekeeper/src/main/etc/cron.d/HouseKeeping
# Need to give proper time
# change the as-vip - required IP.
0 0 * * * root curl -d "querySql=DELETE FROM events WHERE created <= DATE_SUB(NOW(), INTERVAL 6 MONTH)" http://as-vip:8803/housekeeping/execute 
